# target/linux/sunxi/base-files/etc/hotplug.d/usb/10-modem_switch
#!/bin/sh
# Copyright (C) 2019 Aleksander Morgado <aleksander@aleksander.es>

# We need to process only full USB device removal events, we don't
# want to process specific interface removal events.
[ "$ACTION" = remove ] || exit


if [ "$(cat /sys/bus/usb/devices/1-1.1/idVendor)" = "12d1" ] || [ "$(cat /sys/bus/usb/devices/1-1.1/idVendor)" = "03f0" ]; then
if [ "$(cat /sys/bus/usb/devices/1-1.1/idProduct)" = "581d" ] || [ "$(cat /sys/bus/usb/devices/1-1.1/idProduct)" = "1573" ]; then
	echo 3 > /sys/bus/usb/devices/1-1.1/bConfigurationValue
fi
fi

if [ "$(cat /sys/bus/usb/devices/1-1.2/idVendor)" = "12d1" ] || [ "$(cat /sys/bus/usb/devices/1-1.2/idVendor)" = "03f0" ]; then
if [ "$(cat /sys/bus/usb/devices/1-1.2/idProduct)" = "581d" ] || [ "$(cat /sys/bus/usb/devices/1-1.2/idProduct)" = "1573" ]; then
	echo 3 > /sys/bus/usb/devices/1-1.2/bConfigurationValue
fi
fi

if [ "$(cat /sys/bus/usb/devices/1-1.3/idVendor)" = "12d1" ] || [ "$(cat /sys/bus/usb/devices/1-1.3/idVendor)" = "03f0" ]; then
if [ "$(cat /sys/bus/usb/devices/1-1.3/idProduct)" = "581d" ] || [ "$(cat /sys/bus/usb/devices/1-1.3/idProduct)" = "1573" ]; then
	echo 3 > /sys/bus/usb/devices/1-1.3/bConfigurationValue
fi
fi

if [ "$(cat /sys/bus/usb/devices/1-1.4/idVendor)" = "12d1" ] || [ "$(cat /sys/bus/usb/devices/1-1.4/idVendor)" = "03f0" ]; then
if [ "$(cat /sys/bus/usb/devices/1-1.4/idProduct)" = "581d" ] || [ "$(cat /sys/bus/usb/devices/1-1.4/idProduct)" = "1573" ]; then
	echo 3 > /sys/bus/usb/devices/1-1.4/bConfigurationValue
fi
fi

if [ "$(cat /sys/bus/usb/devices/3-1.1/idVendor)" = "12d1" ] || [ "$(cat /sys/bus/usb/devices/3-1.1/idVendor)" = "03f0" ]; then
if [ "$(cat /sys/bus/usb/devices/3-1.1/idProduct)" = "581d" ] || [ "$(cat /sys/bus/usb/devices/3-1.1/idProduct)" = "1573" ]; then
	echo 3 > /sys/bus/usb/devices/3-1.1/bConfigurationValue
fi
fi

if [ "$(cat /sys/bus/usb/devices/3-1.2/idVendor)" = "12d1" ] || [ "$(cat /sys/bus/usb/devices/3-1.2/idVendor)" = "03f0" ]; then
if [ "$(cat /sys/bus/usb/devices/3-1.2/idProduct)" = "581d" ] || [ "$(cat /sys/bus/usb/devices/3-1.2/idProduct)" = "1573" ]; then
	echo 3 > /sys/bus/usb/devices/1-1.2/bConfigurationValue
fi
fi

if [ "$(cat /sys/bus/usb/devices/3-1.3/idVendor)" = "12d1" ] || [ "$(cat /sys/bus/usb/devices/3-1.3/idVendor)" = "03f0" ]; then
if [ "$(cat /sys/bus/usb/devices/3-1.3/idProduct)" = "581d" ] || [ "$(cat /sys/bus/usb/devices/3-1.3/idProduct)" = "1573" ]; then
	echo 3 > /sys/bus/usb/devices/3-1.3/bConfigurationValue
fi
fi

if [ "$(cat /sys/bus/usb/devices/3-1.4/idVendor)" = "12d1" ] || [ "$(cat /sys/bus/usb/devices/3-1.4/idVendor)" = "03f0" ]; then
if [ "$(cat /sys/bus/usb/devices/3-1.4/idProduct)" = "581d" ] || [ "$(cat /sys/bus/usb/devices/3-1.4/idProduct)" = "1573" ]; then
	echo 3 > /sys/bus/usb/devices/3-1.4/bConfigurationValue
fi
fi


